# docker/crontab
10 * * * * www php think VpsTest >> /var/log/cron.log 2>&1